# ansible/templates/deploy_fds/upstart/fds-om.conf.j2
{% extends 'fds_upstart_defaults.j2' %}
{% block content %}
##
# Orchestration Manager upstart conf file
#
description "Orchestration Manager"
version     "2015.07.03-SNAPSHOT"

start on runlevel [2345]
stop on runlevel [!2345]

env fds_root={{ fds_root }}
env ENABLED=false

limit nofile 400000 400000
limit core unlimited unlimited

script
  export fds_root
  chdir ${fds_root}/bin

  # don't start in containers
  /bin/running-in-container && exit 0

  [ -f ${fds_root}/etc/omenv ] || exit 0
  . ${fds_root}/etc/omenv
  [ "${ENABLED}" = "true" ] || exit 0

  echo '\n\n\n=====================================================' > ${fds_root}/bin/om.out
  echo 'UPTIME:\t'`date` >> ${fds_root}/bin/om.out
  echo 'ENV   :\t'`set` >> ${fds_root}/bin/om.out
  echo '=====================================================\n\n\n' >> ${fds_root}/bin/om.out

  exec ${JAVA_HOME}/bin/java ${JAVA_OPTS} -cp ${CLASSPATH} ${MAIN_CLASS} --foreground --fds-root=${fds_root} >> ${fds_root}/bin/om.out 2>&1
end script

{% endblock %}
